# pamsignal.service
[Unit]
Description=PAMSignal - Real-time PAM login monitor
Documentation=https://github.com/anhtuank7c/pamsignal
After=systemd-journald.service

[Service]
Type=simple
ExecStart=/usr/local/bin/pamsignal --foreground
User=pamsignal
Group=pamsignal
SupplementaryGroups=systemd-journal
Restart=on-failure
RestartSec=5

# Runtime directory for PID file (/run/pamsignal/)
RuntimeDirectory=pamsignal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
PrivateDevices=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target
